# deploy/viewer-360.supervisor.conf
[program:viewer-360]
command=python3 -m http.server 4301 --directory /home/jose_luu/360-library-viewer
directory=/home/jose_luu/360-library-viewer
autostart=true
autorestart=true
stderr_logfile=/var/log/viewer-360.err.log
stdout_logfile=/var/log/viewer-360.out.log
user=jose_luu
